# docker/root/etc/services.d/filelia/run
#!/usr/bin/with-contenv sh

/app/start.sh
